# systemd/user/batsignal.service.d/options.conf
[Service]
ExecStart=
ExecStart=batsignal -f 97 -p -m 30
